# nexos-live/config/includes.chroot/etc/profile.d/nexos-live.sh
#!/bin/sh
export NEXOS_LIVE=1
export NEXOS_UI=/usr/share/nexos-ui/index.html
